# roles/vault_servers/files/systemd-vault-unseal.service
# {{ ansible_managed }}

[Unit]
Description=Vault Unseal
After=vault.service
Requires=vault.service
PartOf=vault.service

[Service]
Type=oneshot
User=vault
Group=vault
ExecStartPre=/bin/sleep 10
ExecStart=/bin/sh -c '/opt/vault/vault-unseal.sh'
RemainAfterExit=false

[Install]
WantedBy=multi-user.target vault.service
